# roles/scroller/files/scroller.cron
*/5 * * * * root /usr/local/sbin/update-scroller
