# roles/logrotate/files/logrotate.cron
 0 0 * * * root logrotate /etc/logrotate.conf
