# roles/kiosk/files/scroller.cron
0 * * * * kiosk /home/kiosk/update-scroller.sh
